# dotfiles/services/change_wallpaper.timer
[Unit]
Description=wallpaper changer timer

[Timer]
OnBootSec=0min
AccuracySec=1us
OnUnitActiveSec=15sec

[Install]
WantedBy=multi-user.target
